# ckan/setup/cron.conf
[program:cron]
command=supercronic /etc/crontabs/ckan
autostart=true
autorestart=true
