# ckan/setup/harvest.conf
; ===============================
; ckan harvester
; ===============================

[program:ckan_gather_consumer]

command=ckan --config=/srv/app/ckan.ini harvester gather-consumer

; user that owns virtual environment.
user=ckan

numprocs=1
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stdout
autostart=true
autorestart=true
startsecs=10

[program:ckan_fetch_consumer]

command=ckan --config=/srv/app/ckan.ini harvester fetch-consumer

; user that owns virtual environment.
user=ckan

numprocs=1
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stdout
autostart=true
autorestart=true
startsecs=10
